# mac/pkg_scripts/postinstall
#!/bin/bash
# mac/pkg_scripts/postinstall
#
# Runs as root AFTER files are copied.
# Sources the temp config written by preinstall, calls the config writer binary,
# registers the launchd service, then removes the temp file (it contains secrets).

set -e

# ── Logging setup ─────────────────────────────────────────────────────────────
LOGFILE="/var/log/genreport_install.log"
exec > >(tee -a "$LOGFILE") 2>&1  # Redirect all stdout+stderr to log file

log() { echo "[$(date '+%Y-%m-%d %H:%M:%S')] $*"; }
die() { log "FATAL: $*"; exit 1; }

log "========== Genreport postinstall started =========="
log "Running as: $(whoami)"
log "macOS version: $(sw_vers -productVersion)"

# shellcheck source=/dev/null
log "Sourcing install config from /tmp/genreport_install_config..."
[ -f /tmp/genreport_install_config ] || die "Config file not found at /tmp/genreport_install_config"
source /tmp/genreport_install_config
log "Config sourced successfully."

INSTALL_DIR="/Applications/Genreport"

log "Checking install directory: $INSTALL_DIR"
[ -d "$INSTALL_DIR" ] || die "Install directory $INSTALL_DIR does not exist"

# ── Install PostgreSQL (if requested by preinstall) ───────────────────────────
if [ "$POSTGRES_INSTALL" = "yes" ]; then
  log "=== PostgreSQL installation requested ==="
  PG_INSTALLED=false
  PG_BIN_DIR=""

  # Identify the actual console user (postinstall runs as root; brew must not)
  CONSOLE_USER=$(stat -f "%Su" /dev/console 2>/dev/null || echo "")
  log "Console user: ${CONSOLE_USER:-unknown}"

  # ── Option 1: Homebrew ──────────────────────────────────────────────────────
  BREW_BIN=""
  for b in "/opt/homebrew/bin/brew" "/usr/local/bin/brew"; do
    [ -x "$b" ] && BREW_BIN="$b" && break
  done

  if [ -n "$BREW_BIN" ] && [ -n "$CONSOLE_USER" ] && [ "$CONSOLE_USER" != "root" ]; then
    log "Homebrew found at $BREW_BIN — installing postgresql@17 as $CONSOLE_USER..."
    if sudo -u "$CONSOLE_USER" "$BREW_BIN" install postgresql@17; then
      log "PostgreSQL installed via Homebrew."
      sudo -u "$CONSOLE_USER" "$BREW_BIN" services start postgresql@17 \
        && log "postgresql@17 service started." \
        || log "Warning: brew services start returned non-zero (may already be running)"
      PG_BIN_DIR=$(sudo -u "$CONSOLE_USER" "$BREW_BIN" --prefix postgresql@17 2>/dev/null)/bin
      PG_INSTALLED=true
    else
      log "Homebrew install failed — falling back to EDB installer."
    fi
  else
    [ -z "$BREW_BIN" ] && log "Homebrew not found — skipping Homebrew path."
    [ "$CONSOLE_USER" = "root" ] && log "Console user is root — skipping Homebrew path."
  fi

  # ── Option 2: Bundled EDB installer ─────────────────────────────────────────
  if [ "$PG_INSTALLED" = "false" ]; then
    EDB_INSTALLER="$INSTALL_DIR/postgresql-17-installer.run"
    log "Looking for bundled EDB installer at: $EDB_INSTALLER"
    [ -f "$EDB_INSTALLER" ] || die "EDB PostgreSQL installer not found in $INSTALL_DIR. Reinstall the package."

    log "Running EDB PostgreSQL 17 installer (headless)..."
    chmod +x "$EDB_INSTALLER"
    "$EDB_INSTALLER" \
      --mode unattended \
      --superpassword "$DB_PASS" \
      --serverport "$DB_PORT" \
      --prefix /Library/PostgreSQL/17 \
      --datadir /Library/PostgreSQL/17/data \
      --enable-components server,commandlinetools \
      && log "PostgreSQL installed via EDB." \
      || die "EDB PostgreSQL installation failed."

    PG_BIN_DIR="/Library/PostgreSQL/17/bin"
    PG_INSTALLED=true
  fi

  export PATH="$PG_BIN_DIR:$PATH"
  log "pg bin dir: $PG_BIN_DIR"

  # ── Wait for server to become ready (max 30s) ───────────────────────────────
  log "Waiting for PostgreSQL to become ready..."
  PG_READY=false
  for i in $(seq 1 30); do
    if pg_isready -h "$DB_HOST" -p "$DB_PORT" &>/dev/null 2>&1; then
      log "PostgreSQL ready after ${i}s."
      PG_READY=true
      break
    fi
    sleep 1
  done
  [ "$PG_READY" = "true" ] || die "PostgreSQL did not become ready within 30 seconds."

  # ── Ensure 'postgres' superuser role exists with the configured password ─────
  # Homebrew installs create a cluster owned by the macOS user (not 'postgres').
  # We create/alter the postgres role so Genreport can always use it.
  log "Configuring postgres superuser role and password..."
  SETUP_SQL="DO \$\$ BEGIN
    IF NOT EXISTS (SELECT FROM pg_roles WHERE rolname = 'postgres') THEN
      CREATE ROLE postgres SUPERUSER LOGIN PASSWORD '$DB_PASS';
    ELSE
      ALTER ROLE postgres SUPERUSER LOGIN PASSWORD '$DB_PASS';
    END IF;
  END \$\$;"

  # Attempt as postgres user first (EDB path), then as console user (Homebrew path)
  if ! ( PGPASSWORD="$DB_PASS" psql -h "$DB_HOST" -p "$DB_PORT" -U postgres -d postgres -c "$SETUP_SQL" &>/dev/null 2>&1 ); then
    if [ -n "$CONSOLE_USER" ] && [ "$CONSOLE_USER" != "root" ]; then
      sudo -u "$CONSOLE_USER" psql -d postgres -c "$SETUP_SQL" \
        && log "postgres role configured via console user." \
        || log "Warning: could not configure postgres role — manual setup may be needed."
    fi
  else
    log "postgres role configured successfully."
  fi
fi

# ── Run the config writer ─────────────────────────────────────────────────────

log "Building configwriter args..."

CW_ARGS=(
  --installdir  "$INSTALL_DIR"
  --appport     "$APP_PORT"
  --dbhost      "$DB_HOST"
  --dbport      "$DB_PORT"
  --dbname      "$DB_NAME"
  --dbuser      "$DB_USER"
  --dbpassword  "$DB_PASS"
  --smtphost    "$SMTP_HOST"
  --smtpport    "$SMTP_PORT"
  --smtpuser    "$SMTP_USER"
  --smtppass    "$SMTP_PASS"
  --smtpfrom    "$SMTP_FROM"
  --r2accountid "$R2_ACCOUNT"
  --r2bucket    "$R2_BUCKET"
  --r2accesskey "$R2_KEY"
  --r2secretkey "$R2_SECRET"
  --r2publicurl "$R2_URL"
)

if [ -n "$JWT_SECRET" ] && [ "$JWT_SECRET" != "(auto-generate)" ]; then
  CW_ARGS+=(--jwtsecret "$JWT_SECRET")
  log "JWT secret: using custom value"
else
  log "JWT secret: will auto-generate"
fi

if [ -n "$ENCRYPTION_KEY" ] && [ "$ENCRYPTION_KEY" != "(auto-generate)" ]; then
  CW_ARGS+=(--encryptionkey "$ENCRYPTION_KEY")
  log "Encryption key: using custom value"
else
  log "Encryption key: will auto-generate"
fi

log "Checking configwriter binary..."
[ -x "$INSTALL_DIR/configwriter" ] || die "configwriter not found or not executable at $INSTALL_DIR/configwriter"

log "Running configwriter..."
"$INSTALL_DIR/configwriter" "${CW_ARGS[@]}" && log "configwriter completed successfully." || die "configwriter failed with exit code $?"

# ── Install and Configure Ollama ──────────────────────────────────────────────
log "OLLAMA_INSTALL=$OLLAMA_INSTALL"

if [ "$OLLAMA_INSTALL" = "yes" ]; then
  log "Downloading and installing Ollama..."
  curl -fsSL https://ollama.com/install.sh | sh && log "Ollama installed successfully." || die "Ollama install script failed"
fi

if [ "$OLLAMA_INSTALL" = "yes" ] || [ "$OLLAMA_INSTALL" = "already_installed" ]; then
  OLLAMA_PID=""
  log "Checking if Ollama server is already running..."
  if ! curl -s http://localhost:11434/ > /dev/null; then
    log "Ollama not running — starting temporary server..."
    ollama serve >/dev/null 2>&1 &
    OLLAMA_PID=$!
    log "Ollama server PID: $OLLAMA_PID — waiting 3s..."
    sleep 3
  else
    log "Ollama server already running."
  fi

  log "Pulling nomic-embed-text model..."
  ollama pull nomic-embed-text && log "Model pulled successfully." || die "ollama pull failed"

  if [ -n "$OLLAMA_PID" ]; then
    log "Stopping temporary Ollama server (PID $OLLAMA_PID)..."
    kill "$OLLAMA_PID"
    log "Ollama server stopped."
  fi
fi

# ── Register launchd service ──────────────────────────────────────────────────
PLIST_SRC="$INSTALL_DIR/com.genreport.launcher.plist"
PLIST_DST="/Library/LaunchDaemons/com.genreport.launcher.plist"

log "Checking plist file: $PLIST_SRC"
[ -f "$PLIST_SRC" ] || die "Plist file not found: $PLIST_SRC"

log "Copying plist to LaunchDaemons..."
cp "$PLIST_SRC" "$PLIST_DST" && log "Plist copied." || die "Failed to copy plist"

log "Loading launchd service..."
launchctl load "$PLIST_DST" && log "launchd service loaded successfully." || die "launchctl load failed"

# ── Clean up ──────────────────────────────────────────────────────────────────
log "Removing temp config file..."
rm -f /tmp/genreport_install_config
log "Cleanup done."

log "========== Genreport postinstall completed successfully =========="
exit 0